# deploy/supervisor_include/simple_demo.conf
[program:app_deployment]
directory = /Users/hz/Desktop/hz/Flask/flask-learn/
command=/Users/hz/Desktop/envs/flask-learn/bin/python app_deployment.py





